# conf/freshrss.cron
*/10 * * * * www-data /usr/bin/php __FINALPATH__/app/actualize_script.php >/tmp/FreshRSS.log 2>&1
