# conf/freshrss.cron
*/10 * * * * __APP__ /usr/bin/php7.0 __FINALPATH__/app/actualize_script.php >/tmp/FreshRSS.log 2>&1
