# conf/freshrss.cron
*/10 * * * * __APP__ /usr/bin/php7.3 __FINALPATH__/app/actualize_script.php >/tmp/__APP__.log 2>&1
