# conf/freshrss.cron
MAILTO="root"
*/10 * * * * __APP__ /usr/bin/php__PHPVERSION__ __INSTALL_DIR__/app/actualize_script.php >/var/log/__APP__/__APP__.log 2>&1
